# Deploy/update/servercenter-agent-update.timer
[Unit]
Description=Periodic ServerCenter agent self-update

[Timer]
# DEV cadence: pull frequently while iterating (daily is the eventual default). Persistent catches
# up a run missed while the node was off.
# EARMARK (2026-07-10): make this interval a controller-managed setting (global / per-node update
# cadence) instead of a baked-in value, so it can be tweaked without reinstalling.
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
